# vps-latency-probe.service
[Unit]
Description=VPS Latency Probe Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vps-latency-probe
ExecStart=/usr/bin/python3 /opt/vps-latency-probe/app.py
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
